# supervisord.conf
[supervisord]
nodaemon=true

[program:multiplexer]
command=/usr/local/bin/qgsd-multiplexer

[program:qgsd]
command=/opt/intel/tdx-qgs/qgs --no-daemon