# tool/supervisord.conf
[supervisord]
nodaemon=true

[program:multiplexer]
command=/usr/local/bin/qgsd-multiplexer

[program:qgsd]
command=/usr/local/bin/qgs -uds_path=/var/run/tdx-qgs/qgs.socket --no-daemon